# packages/spin-hmi/spin_hmi-0.9.2.tar.gz/spin_hmi-0.9.2/etc/spin.service
[Unit]
Description=Spin daemon
After=entangle.target frappy.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/spin
ExecStart=/usr/bin/spin
Restart=on-abnormal
RestartForceExitStatus=3

[Install]
WantedBy=multi-user.target
