# packages/gmailnudge/gmailnudge-2.1.3.tar.gz/gmailnudge-2.1.3/src/gmailnudge/deployment_files/gmailnudge.service
[Unit]
Description=Nudge GMail to force frequent new email checks from my server
After=network-online.target

[Service]
ExecStart=/path-to-venv/bin/sndemail

[Install]
WantedBy=multi-user.target
