# packages/mcp-generate-mind/mcp_generate_mind-0.1.0.tar.gz/mcp_generate_mind-0.1.0/supervisord.conf
[unix_http_server]
file=/app/run/generate_mind.sock   ; (the path to the socket file)

[supervisord]
nodaemon=true
logfile=/app/log/supervisord.log
pidfile=/app/run/supervisord.pid
user=root
directory=/app/log
childlogdir=/app/log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:////app/run/generate_mind.sock

[program:generate_mind]
command=uv run python ./src/mcp-generate-mind/start_server.py
directory=/app
redirect_stderr=true
stdout_logfile=/app/log/supervisor
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=NONE
