# packages/tessdb-server/tessdb_server-5.0.9.tar.gz/tessdb_server-5.0.9/files/tdbalarm.timer
# Ensure to do the following:
# sudo systemctl daemon-reload
# sudo systemctl enable tdbalarm.timer
# sudo systemctl start tdbalarm.timer


[Unit]
Description=Runs every 20 minutes

[Timer]
OnCalendar=*-*-* *:0/20:00

[Install]
WantedBy=timers.target
